# clipper.cron
0 5 * * * root /usr/local/clipper/bin/clipper
